# templates/factorio.service.erb
[Unit]
Description=Factorio Headless Server
After=network.target

[Service]
Type=simple
User=factorio
WorkingDirectory=/home/factorio
ExecStart=/home/factorio/factorio/bin/x64/factorio --start-server /home/factorio/factorio/saves/world.zip
Restart=on-abort

[Install]
WantedBy=multi-user.target